# clawcode-bridge.service
[Unit]
Description=Clawcode Telegram bridge
After=opencode-server.service
Requires=opencode-server.service

[Service]
Type=simple
WorkingDirectory={{WORKDIR}}
EnvironmentFile={{WORKDIR}}/.env
ExecStartPre=/usr/bin/curl -sf --retry 30 --retry-delay 1 --retry-connrefused http://127.0.0.1:4096/session
ExecStart=/usr/bin/bun run src/main.ts
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
