# opencode-server.service
[Unit]
Description=OpenCode server
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKDIR}}
ExecStart=/usr/bin/opencode serve --port 4096
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
